# openproject/vendor/pkgr/scaling/upstart/openproject-console-PROCESS_NUM.conf
start on starting openproject-console
stop on stopping openproject-console
respawn

env PORT=PORT_NUM

exec openproject run console >> /var/log/openproject/console-PROCESS_NUM.log 2>&1
