# openproject/vendor/pkgr/scaling/upstart/openproject-check-PROCESS_NUM.conf
start on starting openproject-check
stop on stopping openproject-check
respawn

env PORT=PORT_NUM

exec openproject run check >> /var/log/openproject/check-PROCESS_NUM.log 2>&1
