# openproject/vendor/pkgr/scaling/upstart/openproject-web-PROCESS_NUM.conf
start on starting openproject-web
stop on stopping openproject-web
respawn

env PORT=PORT_NUM

exec openproject run web >> /var/log/openproject/web-PROCESS_NUM.log 2>&1
