# openproject/vendor/pkgr/scaling/upstart/openproject-backup.conf
start on starting openproject
stop on stopping openproject
