# openproject/vendor/pkgr/scaling/upstart/openproject-backup-PROCESS_NUM.conf
start on starting openproject-backup
stop on stopping openproject-backup
respawn

env PORT=PORT_NUM

exec openproject run backup >> /var/log/openproject/backup-PROCESS_NUM.log 2>&1
