# openproject/vendor/pkgr/scaling/upstart/openproject-worker-PROCESS_NUM.conf
start on starting openproject-worker
stop on stopping openproject-worker
respawn

env PORT=PORT_NUM

exec openproject run worker >> /var/log/openproject/worker-PROCESS_NUM.log 2>&1
