# openproject/vendor/pkgr/scaling/upstart/openproject-rake-PROCESS_NUM.conf
start on starting openproject-rake
stop on stopping openproject-rake
respawn

env PORT=PORT_NUM

exec openproject run rake >> /var/log/openproject/rake-PROCESS_NUM.log 2>&1
